# roles/acme-tiny/files/acme-tiny.timer
[Unit]
Description=ACME-Tiny (timer)

[Timer]
OnCalendar=monthly

[Install]
WantedBy=timers.target
